# tools/openclaw-boot-health/openclaw-boot-health-check.service
[Unit]
Description=OpenClaw boot health check — auto-detect + heal after reboot
After=openclaw-gateway.service
Requires=openclaw-gateway.service
PartOf=openclaw-gateway.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/scripts/openclaw-boot-health-check.py
TimeoutStartSec=60
RemainAfterExit=no

[Install]
WantedBy=default.target
