# tools/openclaw-session-backup/openclaw-session-backup.timer
[Unit]
Description=OpenClaw Session Backup Timer — 每 10 分钟备份会话状态

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
