# tools/embed-sidecar/openclaw-embed-sidecar.service
# 适用机器：公司（Linux）
# 系统 / OS：Linux
# 用途：本地向量模型常驻 HTTP 服务（embed-sidecar），消除每次搜索 10s 模型加载

[Unit]
Description=OpenClaw embed sidecar (vector model)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.openclaw/workspace
ExecStart=%h/.local/share/openclaw-embed-venv311/bin/python3 %h/.openclaw/workspace/scripts/embed-sidecar.py --host 127.0.0.1 --port 18792
Restart=on-failure
RestartSec=2
TimeoutStartSec=120
# 模型加载需要 ~12s，给足时间

[Install]
WantedBy=default.target
