# tools/openclaw-watchers/openclaw-task-scheduler.service
[Unit]
Description=Periodic task scheduler check — auto supervisor, stall notify, zombie kill
After=network-online.target
Wants=network-online.target

[Service]
Environment=PATH=/home/missyouangeled/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
Type=oneshot
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/scripts/openclaw-task-scheduler.py --print-json
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-task-scheduler
