# tools/openclaw-session-watcher/openclaw-session-watcher.service
[Unit]
Description=OpenClaw Session Size Watcher (cron fallback)
After=default.target

[Service]
Type=oneshot
ExecStart=python3 %h/.openclaw/workspace/scripts/openclaw-session-size-watcher.py --systemd
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
