# mark42-pkg/mark42/systemd/mark42-watchdog.timer
[Unit]
Description=Mark42 Watchdog Timer

[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=10

[Install]
WantedBy=default.target
