# mark42-pkg/mark42/systemd/mark42-bootstrap.service.tmpl
[Unit]
Description=Mark42 Bootstrap - 开机自动挂载守护 Loop
After=openclaw-gateway.service
Wants=openclaw-gateway.service

[Service]
Type=oneshot
Environment=MARK42_WORKSPACE=__MARK42_WORKSPACE__
Environment=MARK42_CLI=__MARK42_BIN__
Environment=MARK42_PYTHON_BIN=__MARK42_PYTHON__
Environment=XDG_STATE_HOME=__MARK42_XDG_STATE__
Environment=MARK42_STATE_DIR=__MARK42_STATE_DIR__
Environment=MARK42_LOG_DIR=__MARK42_LOG_DIR__
Environment=MARK42_SCRATCH=__MARK42_SCRATCH__
ExecStart=__MARK42_PYTHON__ __MARK42_BIN__ engine register --all
StandardOutput=journal
StandardError=journal
RemainAfterExit=yes

[Install]
WantedBy=default.target
