# tools/mark42-engine-daemon/mark42-engine-daemon.service
# TEMPLATE FILE — 需先渲染 __MARK42_* 占位符，再安装到 ~/.config/systemd/user/
[Unit]
Description=Mark42 Engine Daemon — 循环引擎守护进程
After=openclaw-gateway.service
Wants=openclaw-gateway.service
Requires=mark42-bootstrap.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
Environment=MARK42_WORKSPACE=__MARK42_WORKSPACE__
Environment=MARK42_CLI=__MARK42_CLI__
Environment=MARK42_PYTHON_BIN=__MARK42_PYTHON__
Environment=XDG_STATE_HOME=__XDG_STATE_HOME__
Environment=MARK42_STATE_DIR=__MARK42_STATE_DIR__
Environment=MARK42_LOG_DIR=__MARK42_LOG_DIR__
Environment=MARK42_SCRATCH=__MARK42_OPENCLAW_SCRATCH_ROOT__
ExecStart=__MARK42_PYTHON__ -u __MARK42_CLI__ engine --daemon --interval 30
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=append:__MARK42_LOG_DIR__/engine-daemon.log
StandardError=append:__MARK42_LOG_DIR__/engine-daemon.log
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=default.target
