# tools/openclaw-local-health/openclaw-local-health-watch.timer
[Unit]
Description=Periodic OpenClaw local health diagnosis

[Timer]
Unit=openclaw-local-health-watch.service
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
