# tools/openclaw-supervisor/openclaw-supervisor-watch.timer
[Unit]
Description=Periodic OpenClaw supervisor status refresh

[Timer]
Unit=openclaw-supervisor-watch.service
OnBootSec=1min
OnUnitActiveSec=30s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
