# tools/openclaw-systemd/openclaw-resume-watch.timer
[Unit]
Description=Watch for sleep/resume gaps and recover OpenClaw

[Timer]
Unit=openclaw-resume-watch.service
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
